# src/celery-worker.service
[Unit]
Description=Celery Worker Service
After=network.target redis.service

[Service]
User=magellan
Group=magellan
WorkingDirectory=/home/magellan/Projects/web200ok
Environment="PATH=/home/magellan/Projects/web200ok/.venv/bin"
ExecStart=/home/magellan/Projects/web200ok/.venv/bin/celery -A config worker -l info -c 4
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target